# board/car-player/rootfs_overlay/etc/udev/rules.d/10-automount-usb-disk.rules
SUBSYSTEMS=="usb", ACTION=="add", KERNEL=="sd*", ATTR{partition}=="?", RUN+="/root/automount.sh add %k"
SUBSYSTEMS=="block", ACTION=="remove", KERNEL=="sd*", ENV{DEVTYPE}=="partition", RUN+="/root/automount.sh remove %k"
